# crontab
00 * * * * ~/web_scripts/startup.mit.edu/New-Website/pull-master.sh
